# com.user.qwenassistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.qwenassistant</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/admin/Developer/stt-tts-qwen/start_assistant.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/admin/Developer/stt-tts-qwen/assistant_output.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/admin/Developer/stt-tts-qwen/assistant_error.log</string>
</dict>
</plist>